# aiidalab/etc/service/jupyter-notebook/run
#!/bin/bash
set -em

/opt/start-notebook.sh

